# deploy/systemd/oldclaw-subagent.service
[Unit]
Description=OldClaw SubAgent Service
After=network.target

[Service]
ExecStart=/usr/local/bin/oldclaw-subagent
Restart=always

[Install]
WantedBy=multi-user.target
