# deploy/systemd/opsclaw-manager.service
[Unit]
Description=OpsClaw Manager Service
After=network.target

[Service]
ExecStart=/usr/local/bin/opsclaw-manager
Restart=always

[Install]
WantedBy=multi-user.target
