# deploy/systemd/opsclaw-subagent.service
[Unit]
Description=OpsClaw SubAgent Service
After=network.target

[Service]
ExecStart=/usr/local/bin/opsclaw-subagent
Restart=always

[Install]
WantedBy=multi-user.target
