# deploy/systemd/opsclaw-subagent.service
[Unit]
Description=OpsClaw SubAgent Runtime
After=network.target

[Service]
Type=simple
User=opsclaw
Environment=PYTHONPATH=/home/opsclaw/opsclaw
ExecStart=/home/opsclaw/opsclaw/.venv/bin/uvicorn apps.subagent-runtime.src.main:app --host 0.0.0.0 --port 8002
WorkingDirectory=/home/opsclaw/opsclaw
EnvironmentFile=/home/opsclaw/opsclaw/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
