# deploy/systemd/opsclaw-master.service
[Unit]
Description=OpsClaw Master Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=opsclaw
Environment=PYTHONPATH=/home/opsclaw/opsclaw
ExecStart=/home/opsclaw/opsclaw/.venv/bin/uvicorn apps.master-service.src.main:app --host 0.0.0.0 --port 8001 --log-level info
WorkingDirectory=/home/opsclaw/opsclaw
EnvironmentFile=/home/opsclaw/opsclaw/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
